# deploy/supervisor/codflow-worker.conf
[program:codflow-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/codflow/artisan queue:work redis --sleep=3 --tries=3 --timeout=120
directory=/var/www/codflow
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/codflow/storage/logs/worker.log
stopwaitsecs=130
